# provisioning/ansible/roles/tomcat/templates/tomcat.service.j2
[Service]
ExecStart={{ tomcat_path }}/bin/catalina.sh run

User={{ tomcat_user }}
Group={{ tomcat_group }}
