# provisioning/ansible/roles/mongodb/sharding/templates/mongosh.service.j2
[Service]
ExecStart={{ mongodb_path }}/bin/mongod --shardsvr --replSet rs1 \
    --config {{ mongodb_shard_config }}

User={{ mongodb_user }}
Group={{ mongodb_group }}
